# system/bbs.service
[Unit]
Description=Talisman BBS service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
WorkingDirectory=/opt/talisman
ExecStart=/opt/talisman/servo

[Install]
WantedBy=multi-user.target
WantedBy=default.target
