# system/crontab
# m h  dom mon dow   command
0 0 * * * /opt/talisman/bin/bbs-logrotate.sh

